# deploy/aws/dregg-gateway.service
[Unit]
Description=Dregg Gateway Node
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=dregg
Group=dregg
WorkingDirectory=/opt/dregg
EnvironmentFile=-/etc/dregg/node.env
ExecStart=/opt/dregg/target/release/dregg-node run \
    --data-dir /opt/dregg-data \
    --port 8420 \
    --bind 127.0.0.1 \
    --node-index 0 \
    --federation-size 4 \
    --enable-pruning \
    --enable-faucet \
    --federation-mode solo \
    --consensus blocklace
ExecStartPost=/opt/dregg/deploy/aws/unlock-gateway.sh
Restart=always
RestartSec=5
Environment=RUST_LOG=info
Environment=DREGG_GATEWAY=true

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/dregg-data
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

[Install]
WantedBy=multi-user.target
