# deploy/aws/dregg-discord-bot.service
[Unit]
Description=Dregg Discord Bot
After=network-online.target dregg-gateway.service
Wants=network-online.target
Requires=dregg-gateway.service

[Service]
Type=simple
User=dregg
Group=dregg
WorkingDirectory=/opt/dregg
EnvironmentFile=/etc/dregg/discord-bot.env
ExecStart=/opt/dregg/target/release/dregg-discord-bot
Restart=always
RestartSec=5
Environment=RUST_LOG=info

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/dregg-data /var/lib/dregg-discord-bot
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

[Install]
WantedBy=multi-user.target
