# deploy/aws/pyana-gateway.service
[Unit]
Description=Pyana Gateway Node
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pyana
Group=pyana
WorkingDirectory=/opt/pyana
ExecStart=/opt/pyana/target/release/pyana-node run \
    --data-dir /opt/pyana-data \
    --port 8420 \
    --morpheus \
    --node-index 0 \
    --federation-size 4 \
    --prove-transitions \
    --enable-pruning \
    --enable-faucet
Restart=always
RestartSec=5
Environment=RUST_LOG=info
Environment=PYANA_GATEWAY=true

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/pyana-data
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

[Install]
WantedBy=multi-user.target
